# deploy/systemd/tetris-relay.service
[Unit]
Description=Tetris Multiplayer Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tetris
Group=tetris
WorkingDirectory=/opt/tetris
EnvironmentFile=/etc/tetris/relay.env
ExecStart=/opt/tetris/tetris_relay --port 7777 --meta https://api.example.com
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
# relay 는 디스크에 아무것도 쓰지 않는다 — 전체 읽기 전용.
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
